# cfg/zeus/etc/systemd/system/cleanup.service
[Unit]
Description=Immigrant archive cleanup

[Service]
Type=oneshot
ExecStart=/srv/scripts/cleanup.sh
