# cfg/zeus/etc/systemd/system/transfer.service
[Unit]
Description=Transfer valid archives to internal server
After=cleanup.service

[Service]
Type=oneshot
ExecStart=/srv/scripts/transfer.sh